# dotfiles/etc/udev/rules.d/10-hdmi-change.rules
KERNEL=="card1" SUBSYSTEM=="drm" ACTION=="change" RUN+="/usr/local/bin/reload-hyprctl.sh"
